# personal/settings/systemd/system/nmb.service.d/restart-on-failure.conf
[Service]
Restart=on-failure
RestartSec=5s