# Personal/settings/autostart/discord.desktop
[Desktop Entry]
Encoding=UTF-8
Version=0.9.4
Type=Application
Name=discord
Comment=
Exec=discord
StartupNotify=false
Terminal=false
Hidden=false

